# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:apache]
command=/usr/sbin/httpd -DFOREGROUND

